# deploy/systemd/gymcore-telegram.service
[Unit]
Description=GymCore Telegram bot
After=network.target gymcore-api.service

[Service]
Type=simple
User=gymcore
WorkingDirectory=/opt/gymcore
EnvironmentFile=/opt/gymcore/.env
ExecStart=/opt/gymcore/.venv/bin/python -m app.bots.telegram_main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
